# systemd/orison-web.service.template
[Unit]
Description=ORISON Flask Web UI
After=network.target

[Service]
User={{USER}}
WorkingDirectory={{BASE_DIR}}/web
ExecStart=/usr/bin/python3 {{BASE_DIR}}/web/app.py
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
